# service/youtupi-update.service
[Unit] 
Description=Update Youtupi

[Service] 
Type=oneshot
WorkingDirectory=/home/pi/youtupi/service
ExecStart=/bin/bash update.sh
